# templates/host-resource-management/machine-tier.slice.j2
[Unit]
Description={{ host_resource_management_tier_name | title }} performance domain slice
Before=slices.target

[Slice]
CPUAccounting=yes
CPUWeight={{ host_resource_management_tier.cpu_weight }}
